# slapd/supervisor.conf
[supervisord]
nodaemon=true
[program:slapd]
priority=10
command=/etc/service/slapd/run
autostart=true
autorestart=true
